# device/testing/device-samsung-on7xelte/mount-efs.rc-init
#!/sbin/openrc-run

command="/bin/mount"
command_args="/dev/disk/by-partlabel/EFS /efs/"

depend() {
	after bootmisc
}

start_pre() {
	mkdir -p /efs/
}
