# modem/qcom-diag/diag-router.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name="DIAG"
description="Qualcomm DIAG router"
command="/usr/bin/diag-router"
