# device/main/device-purism-librem5/l5-brcm-bt-enable.initd
#!/sbin/openrc-run
name="Purism Librem 5 Sparklan brcm bluetooth enabler"
supervisor=supervise-daemon
command="/usr/bin/btattach"
command_args="-B /dev/ttymxc3 -P bcm"
