# device/main/device-purism-librem5/l5-lockdown-support.initd
#!/sbin/openrc-run
name="Purism Librem 5 Lockdown Support"

start()
{
	ebegin "Enabling lockdown support"
        /usr/sbin/lockdown-support.sh
	eend 0
}
