# device/main/device-purism-librem5/l5-ship-mode.initd
#!/sbin/openrc-run
name="Purism Librem 5 Ship Mode Enabler"

start()
{
	ebegin "Enabling ship mode"
        l5-poweroff-shipmode poweroff
	eend 0
}
