# device/testing/device-samsung-a6lte/etc-init.d-wlan0-up
#!/sbin/openrc-run

description="Hack to fix wlan0 EINVAL if rfkilled on boot"

start () {
	ifconfig wlan0 up
	ifconfig wlan0 down
}
