# main/postmarketos-bootsplash/kill-pbsplash
#!/sbin/openrc-run
# Copyright 2023 Caleb Connolly
# SPDX-License-Identifier: GPL-3.0-or-later
name="kill-pbsplash"
description="Stop the splash animation"

command="/usr/bin/killall"
command_args="pbsplash"

depend() {
        before display-manager
}
