# device/community/device-purism-librem5/librem5-shipmode.initd
#!/sbin/openrc-run
description="Put the Librem 5 into 'ship mode' so that the battery does not drain when off."
command="/usr/bin/librem5-shipmode"
