# device/testing/pinenote-dbus-service/pinenote.initd
#!/sbin/openrc-run

description="Launch DBus service for driver aspects"
command="/usr/bin/pinenote_dbus_service"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need dbus
}
