# device/community/device-pine64-pinephone/setup-modem.initd
#!/sbin/openrc-run

command="/usr/bin/pinephone_setup-modem"
command_background=true

depend() {
	after eg25
	before ofono
}
