# device/device-nokia-rx51/udev/10-nokia-modem.rules
# Export GPIOs for nokia-modem
DRIVER=="nokia-modem" RUN+="/bin/ln -sf /sys/bus/hsi/devices/n900-modem /dev/cmt"
