# main/shutdown-clear-rtc-wakealarm/clear-rtc-wakealarm.shutdown.initd
#!/sbin/openrc-run

description="Clear RTC wake alarm during shutdown"

start() {
	echo 0 >/sys/class/rtc/rtc0/wakealarm
}
