# device/testing/firmware-samsung-degaswifi/sd8887-init.initd
#!/sbin/openrc-run

description="Initialize sd8887 wlan driver"

depend()
{
	need udev
	before wlan
}

start()
{
	ebegin "Enabling wlan driver"
	echo 1 > /sys/devices/platform/sd8x-rfkill/pwr_ctrl
	echo drv_mode=1 > /proc/mwlan/config
}
