# main/postmarketos-ondev/postmarketos-ondev.initd
#!/sbin/openrc-run

description="postmarketOS on-device installer"

depend()
{
	need localmount dbus
}

start()
{
	ebegin "Running postmarketOS on-device installer"
	ondev-boot &
}
