# main/iio-sensor-proxy/iio-sensor-proxy.initd
#!/sbin/openrc-run

command="/usr/sbin/iio-sensor-proxy"
pidfile="/run/${RC_SVCNAME}.pid"
command_args=""
command_background=true
command_user="root:user"
